# debian/aish-sandbox.service
[Unit]
Description=AI Shell sandbox daemon
After=network.target aish-sandbox.socket
Requires=aish-sandbox.socket
PartOf=aish-sandbox.socket

[Service]
Type=simple
ExecStart=/usr/bin/aish-sandbox

# systemd socket activation passes the listening FD automatically.
# The daemon will use LISTEN_FDS/LISTEN_PID and fd=3.

User=root
Group=root

# Hardening (keep minimal first; can be tightened later)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
PrivateMounts=yes

[Install]
WantedBy=multi-user.target
