# debian/aish.postrm
#!/bin/sh
set -e

case "$1" in
  remove|purge)
    if [ -d /run/systemd/system ] && command -v systemctl >/dev/null 2>&1; then
      systemctl disable --now aish-sandbox.socket >/dev/null 2>&1 || true
      systemctl stop --no-block aish-sandbox.service >/dev/null 2>&1 || true
      systemctl reset-failed aish-sandbox.service >/dev/null 2>&1 || true
      systemctl daemon-reload >/dev/null 2>&1 || true
    fi
  ;;
esac

#DEBHELPER#

exit 0